# srcs/phpmyadmin/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm7]
command=php-fpm7

[program:nginx]
command=nginx -g 'daemon off;'